# configs/config/systemd/user/pull_notes.service
[Unit]
Description=Pull Notes from Git Repo
After=network.target nss-lookup.target

[Service]
ExecStart=%h/tools/pull_notes.sh
Restart=on-failure
RuntimeMaxSec=60
RestartSec=600

[Install]
WantedBy=default.target
